# root/etc/s6-overlay/s6-rc.d/svc-retrostack-emulator/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

# Redirect stdout/stderr to the container's stdout so output appears in docker logs.
exec > /proc/1/fd/1 2>&1
exec /usr/local/bin/retrostack-emulator-run
